# macinstall/pkg/Scripts.d/postinstall
#!/usr/bin/env bash

launchctl unload -w /Library/LaunchDaemons/com.hurpdurp.apps.plist || true
launchctl load -w /Library/LaunchDaemons/com.hurpdurp.apps.plist
